# config/pfal-controller.service
[Unit]
Description=PFAL Controller - Rule-based automation for Plant Factory
After=network.target mosquitto.service influxdb.service
Wants=mosquitto.service influxdb.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pfal-management
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /home/pi/pfal-management/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pfal-controller

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
